# roles/app_install/templates/webapp.service.j2
[Unit]
Description=Simple Flask Web App
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/opt/webapp
ExecStart=/usr/bin/python3 /opt/webapp/app.py
Restart=always

[Install]
WantedBy=multi-user.target